# scripts/english-pilot.service
[Unit]
Description=EnglishPilot external channel daemon
After=network-online.target

[Service]
Type=simple
Environment=ENGLISH_PILOT_HOME=__HOME__/.english-pilot
ExecStart=__NODE_BIN__ __CLI_JS__ run
Restart=on-failure
RestartSec=5
WorkingDirectory=__HOME__

[Install]
WantedBy=default.target
