# scripts/english-pilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.octopusgarage.english-pilot</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__CLI_JS__</string>
    <string>run</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ENGLISH_PILOT_HOME</key>
    <string>__HOME__/.english-pilot</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.english-pilot/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.english-pilot/logs/launchd.err.log</string>
</dict>
</plist>
